# conf/systemd_check.timer.j2
[Unit]
Description=Check {{ app }} backup regularly

[Timer]
OnCalendar={{ check_on_calendar }}

[Install]
WantedBy=timers.target
